# retrobarbershop.conf
[program:api]
command=/home/ubuntu/go/src/retrobarbershop.com/retro/api/api
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/retro.log
stdout_logfile=/var/log/supervisor/retro.log
